# supervisor.conf
[supervisord]
nodaemon=true

[program:minecraft]
command=/start
